# configs/systemd/dist-git-gc.timer
[Unit]
Description=Git Repositories Clean Up Timer

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
