# files/genacls.cron
02 10 * * * root /usr/local/bin/genacls.sh
