# files/genacls.cron
02 10 * * * root /usr/share/dist-git/genacls.sh
