# files/genacls.cron
02 10 * * * root /usr/share/dist-git/dist_git_sync.sh
